# theme-switcher.path
[Unit]
Description=Detect changes in KDE global theme

[Path]
PathChanged=%h/.config/kdeglobals

[Install]
WantedBy=default.target
